# data/configs/autostart/syncthing.desktop
[Desktop Entry]
Type=Application
Name=Syncthing
# Change this to a full path?
Exec=syncthing --no-browser
Hidden=false
