# gentoo/etc/local.d/charging-thresholds.start
#!/bin/sh

echo 50 > /sys/class/power_supply/BAT0/charge_control_end_threshold
echo 40 > /sys/class/power_supply/BAT0/charge_control_start_threshold
